# ansible/roles/weddingbot/templates/weddingbot-burst.timer.j2
# {{ ansible_managed }}
# Burst poll: every minute for a ~20 min window around the 05:00 Copenhagen
# slot drop, Mon-Thu. The bot latches after the first successful alert and skips
# the rest of the window, so this never spams.
#
# The server clock is UTC. To hit 05:00 Copenhagen year-round WITHOUT depending
# on the systemd-252 timezone prefix, we emit BOTH DST placements as UTC ranges:
#   * CEST (summer, UTC+2): 04:50-05:10 local == 02:50-03:10 UTC
#   * CET  (winter, UTC+1): 04:50-05:10 local == 03:50-04:10 UTC
# On any given date only one is the real drop window; the other just polls
# harmlessly (finds nothing, or finds slots and alerts — both fine).
[Unit]
Description=Wedding bot burst poll around 05:00 Copenhagen slot drop (Mon-Thu)

[Timer]
OnCalendar=Mon,Tue,Wed,Thu 02:50..59:00
OnCalendar=Mon,Tue,Wed,Thu 03:00..10:00
OnCalendar=Mon,Tue,Wed,Thu 03:50..59:00
OnCalendar=Mon,Tue,Wed,Thu 04:00..10:00
Persistent=false
Unit=weddingbot.service

[Install]
WantedBy=timers.target
